# lc-rules/windows_sysmon/sysmon_asep_reg_keys_modification.yml
detect:
  log type: wel
  op: and
  rules:
  - case sensitive: false
    op: is
    path: Event/System/EventID
    value: '13'
  - op: or
    rules:
    - case sensitive: false
      op: contains
      path: Event/EventData/TargetObject
      value: \software\Microsoft\Windows\CurrentVersion\Run
    - case sensitive: false
      op: contains
      path: Event/EventData/TargetObject
      value: \software\Microsoft\Windows\CurrentVersion\RunOnce
    - case sensitive: false
      op: contains
      path: Event/EventData/TargetObject
      value: \software\Microsoft\Windows\CurrentVersion\RunOnceEx
    - case sensitive: false
      op: contains
      path: Event/EventData/TargetObject
      value: \software\Microsoft\Windows\CurrentVersion\RunServices
    - case sensitive: false
      op: contains
      path: Event/EventData/TargetObject
      value: \software\Microsoft\Windows\CurrentVersion\RunServicesOnce
    - case sensitive: false
      op: contains
      path: Event/EventData/TargetObject
      value: \software\Microsoft\Windows NT\CurrentVersion\Winlogon\Userinit
    - case sensitive: false
      op: contains
      path: Event/EventData/TargetObject
      value: \software\Microsoft\Windows NT\CurrentVersion\Winlogon\Shell
    - case sensitive: false
      op: contains
      path: Event/EventData/TargetObject
      value: \software\Microsoft\Windows NT\CurrentVersion\Windows
    - case sensitive: false
      op: contains
      path: Event/EventData/TargetObject
      value: \software\Microsoft\Windows\CurrentVersion\Explorer\User Shell Folders
  target: log
respond:
- action: report
  metadata:
    author: Victor Sergeev, Daniil Yugoslavskiy, oscd.community
    description: Detects modification of autostart extensibility point (ASEP) in registry
    falsepositives:
    - Legitimate software automatically (mostly, during installation) sets up autorun
      keys for legitimate reason
    - Legitimate administrator sets up autorun keys for legitimate reason
    level: medium
    references:
    - https://github.com/redcanaryco/atomic-red-team/blob/master/atomics/T1060/T1060.yaml
    tags:
    - attack.persistence
    - attack.t1060
    - attack.t1547.001
  name: Autorun Keys Modification

